# hardware/dock/com.tlc.dock.seed.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tlc.dock.seed</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/lib/tlc/dock/setup-dock.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>600</integer>
  <key>StandardOutPath</key>
  <string>/var/log/tlc-dock-seed.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/tlc-dock-seed.log</string>
</dict>
</plist>
